# scripts/micropanel.service
[Unit]
Description=MicroPanel - Static Hosting Control Panel
After=network.target nginx.service
Wants=nginx.service

[Service]
Type=simple
User=micropanel
Group=micropanel
ExecStart=/usr/bin/micropanel serve
Restart=on-failure
RestartSec=5

# Security
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/micropanel /var/www/panel/sites /etc/nginx/sites-enabled
PrivateTmp=true

[Install]
WantedBy=multi-user.target
